# files/etc/profile.d/packages.brew
tap "homebrew/aliases"
brew "f3"
brew "fisher"
brew "mise"
brew "ncdu"
brew "smartmontools"
brew "wakeonlan""
